# scripts/ots-hub.service
[Unit]
Description=OTS Hub Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ots_hub
EnvironmentFile=/home/ubuntu/ots_hub/.env
ExecStart=/bin/bash /home/ubuntu/ots_hub/scripts/start_hub.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
